# systemd/user/hypraway-batt-update.service
[Unit]
Description=Update Hypraway combined battery status

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update_hypraway_battery_status.sh
